# services/deluged/run
#!/bin/sh
su-exec deluge deluged -d -c /config
