# services/deluge_exporter/run
#!/bin/sh
DELUGE_CONFIG_DIR=/config/ su-exec deluge /opt/deluge_exporter/deluge_exporter.py
