# etc/services.d/openvpn/run
#!/usr/bin/with-contenv sh

# launch openvpn
exec openvpn \
    --config ${OPENVPN_CONFIG_FILE} \
    --auth-user-pass ${CREDENTIALS_FILE} \ # use the created credentials file for authentication \
    --pull-filter ignore "DNS" \ # don't use VPN server's DNS configuration \
    --dhcp-option DNS ${DNS} \ # instead set the one provided \
    --redirect-gateway def1 bypass-dhcp
