# rootfs/etc/services.d/qbittorrent/run
#!/usr/bin/with-contenv sh

CONFIG_FILE=${CONFIG_DIR}/qBittorrent/config/qBittorrent.conf
export QBT_PROFILE=${CONFIG_DIR}

# wait for VPN connection
until ifconfig | grep tun | grep -q "00-00-00-00-00-00-00-00-00-00-00-00-00-00-00-00"; do
    sleep 0.5
done

DNS_SERVER=$(sed -n -e 's/^.*nameserver //p' /etc/resolv.conf)
# drop rule allowing DNS traffic outside VPN tunnel
# but only if DNS server is not docker's builtint resolver (or LAN traffic which is always allowed)
if [ "${DNS_SERVER}" != "127.0.0.11" ]; then
    # NOTE: these rules are dropped before starting qbittorrent
    iptables -D OUTPUT -p udp -d ${DNS_SERVER} --dport 53 -m state --state NEW,ESTABLISHED -j ACCEPT
    iptables -D INPUT -p udp -s ${DNS_SERVER} --sport 53 -m state --state ESTABLISHED -j ACCEPT
    iptables -D OUTPUT -p tcp -d ${DNS_SERVER} --dport 53 -m state --state NEW,ESTABLISHED -j ACCEPT
    iptables -D INPUT -p tcp -s ${DNS_SERVER} --sport 53 -m state --state ESTABLISHED -j ACCEPT
fi

# verify the iptables rule was dropped
if [ $? -ne 0 ]; then
    echo "ERROR: could not drop iptables rule allowing DNS traffic"
    exit 1
fi

exec s6-setuidgid ${PUID}:${PGID} qbittorrent-nox
