# services/nginx/run
#!/bin/sh
rm -f /var/run/nginx.sock
nginx -g 'daemon off;'
