# services/deluge-web/run
#!/bin/sh
su-exec deluge deluge-web -d -c /config
